# config/supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=apache2ctl -DFOREGROUND
killasgroup=true
stopasgroup=true
autorestart=true

[program:sshd]
command=/usr/sbin/sshd -D
